# .exclude/etc/udev/rules.d/95-hotplug-drm.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/bin/sh -c 'cat /sys/class/drm/card0-HDMI-A-2/status > /tmp/.hdmi_status'"
