# .exclude/etc/udev/rules.d/95-powerplug-change.rules
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/bin/sh -c 'cat /sys/class/power_supply/AC/online > /tmp/.ac_status'"
